# deployment/com.whisper.api-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whisper.api-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/micke/whisper-svenska/api_server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/micke/whisper-svenska</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ENABLE_MLX_REALTIME</key>
        <string>true</string>
        <key>ENABLE_GPU</key>
        <string>true</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/whisper-api-server.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/whisper-api-server.log</string>

    <key>ProcessType</key>
    <string>Standard</string>
</dict>
</plist>
